# vm/network/connen/files/etc/init.d/app
#!/bin/sh /etc/rc.common

START=93
STOP=93

start() {
	chmod 755 /root/start.sh
	sh /root/start.sh &
}

stop() {
	echo "App stop" 
}
